# deploy/spark.service
[Unit]
Description=Spark Pod Orchestrator
After=network-online.target nats-server.service
Wants=network-online.target
Requires=nats-server.service

[Service]
Type=simple
EnvironmentFile=/etc/spark/spark.env
ExecStart=/usr/local/bin/spark \
    -nats ${SPARK_NATS_URL} \
    -node-id ${SPARK_NODE_ID} \
    -manifest-dir ${SPARK_MANIFEST_DIR} \
    -gpu-max ${SPARK_GPU_MAX} \
    -heartbeat-interval ${SPARK_HEARTBEAT_INTERVAL} \
    -reconcile-interval ${SPARK_RECONCILE_INTERVAL} \
    -system-reserve-cpu ${SPARK_SYSTEM_RESERVE_CPU} \
    -system-reserve-memory ${SPARK_SYSTEM_RESERVE_MEMORY} \
    -system-reserve-cores ${SPARK_SYSTEM_RESERVE_CORES} \
    -state-db ${SPARK_STATE_DB} \
    -pod-retention ${SPARK_POD_RETENTION} \
    -http-addr ${SPARK_HTTP_ADDR} \
    -shutdown-timeout ${SPARK_SHUTDOWN_TIMEOUT} \
    -reconcile-resources-interval ${SPARK_RECONCILE_RESOURCES_INTERVAL}
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
